# arqs/openoffice.sh
#!/bin/sh

### BEGIN INIT INFO
# Provides:		openoffice.sh
# Required-Start:	$remote_fs $syslog
# Required-Stop:	$remote_fs $syslog
# Should-Start:		$network
# Should-Stop:		$network
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description:	OpenOffice
# Description:		OpenOffice, muito melhor do que o MSOffice.
### END INIT INFO

/usr/bin/soffice --nologo --nofirststartwizard --headless --norestore --invisible "--accept=socket,host=localhost,port=8100,tcpNoDelay=1;urp;" &
